# ipv6tunnel/tunnel.service
[Unit]
Description=IPv6 Tunneling Service
After=network.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/opt/ipv6tunnel/lib/tunnel.sh start
ExecStop=/opt/ipv6tunnel/lib/tunnel.sh stop
ExecReload=/opt/ipv6tunnel/lib/tunnel.sh restart
Restart=on-failure
RestartSec=5
User=root
Group=root

[Install]
WantedBy=multi-user.target
